# debian/usbip-client.postinst
#!/bin/bash

systemctl enable tls_proxy.service
systemctl enable usbip_autoimport.service

systemctl daemon-reload

systemctl start tls_proxy.service
systemctl start usbip_autoimport.service
